# etc/cron.d/motd
* * * * * root /usr/local/sbin/update-motd 
